# deployment/grest.conf
description "uWSGI server instance configured to serve myproject"

start on runlevel [2345]
stop on runlevel [!2345]

exec /var/gator/api/wsgi.py
