# deployment/grest.conf
description "Starts the flask application"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/gator/api/
exec /var/gator/api/wsgi.py
